# templates/supervisord.conf.in
[unix_http_server]
file = ${buildout:directory}/var/supervisor.sock
chmod = 0600
#chown = root:root

[supervisorctl]
serverurl = unix://${buildout:directory}/var/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisord]
logfile = ${buildout:directory}/var/log/supervisord.log
logfile_maxbytes = 5MB
logfile_backups = 10
loglevel = info
pidfile = ${buildout:directory}/var/supervisord.pid ;
childlogdir = ${buildout:directory}/var/log
nodaemon = false
directory = ${buildout:directory}/var/log

[program:nginx]
command = ${nginx:location}/sbin/nginx
environment = LD_LIBRARY_PATH='${nginx:libxml2}/lib:${nginx:libxslt}/lib'
autostart = true
autorestart = true
stdout_logfile = NONE
stderr_logfile = NONE
priority = 12

[program:haproxy]
command = ${buildout:directory}/bin/haproxy -f ${buildout:directory}/etc/haproxy.conf -db
autostart= true
autorestart = true
redirect_stderr = true
priority = 11

# [program:varnish]
# command = $ ####### {varnish:location}/sbin/varnishd -f ${buildout:directory}/etc/varnish.vcl -F
# autostart= true
# autorestart = true
# redirect_stderr = true
# priority = 13

[program:zeo]
command = ${zeo:location}/bin/runzeo 
process_name = zeo
directory = ${zeo:location}
priority = 10
redirect_stderr = false

[program:instance1]
command = ${instance1:location}/bin/runzope 
process_name = instance1
directory = ${instance1:location}
priority = 20
redirect_stderr = false

[program:instance2]
command = ${instance2:location}/bin/runzope 
process_name = instance2
directory = ${instance2:location}
priority = 30
redirect_stderr = false

[program:instance3]
command = ${instance3:location}/bin/runzope 
process_name = instance3
directory = ${instance3:location}
priority = 20
redirect_stderr = false

[program:instance4]
command = ${instance4:location}/bin/runzope 
process_name = instance4
directory = ${instance4:location}
priority = 30
redirect_stderr = false
