# templates/supervisord.conf.in
[supervisord]
umask=002                   ; (process file creation umask;default 022)
logfile=${parts.buildout.directory}/var/log/supervisord.log ; (main log file;default $$CWD/supervisord.log)
logfile_maxbytes=50MB       ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10          ; (num of main logfile rotation backups;default 10)
loglevel=info               ; (logging level;default info; others: debug,warn)
pidfile=${parts.buildout.directory}/var/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false              ; (start in foreground if true;default false)
minfds=1024                 ; (min. avail startup file descriptors;default 1024)
minprocs=200                ; (min. avail process descriptors;default 200)
user=${options.user}
enironment=PATH=/usr/local/sbin:/usr/local/bin:/usr/bin:/usr/sbin:/bin:/sbin:/usr/X11R6/bin
identifier=supervisor

[unix_http_server]
file=${parts.buildout.directory}/var/supervisor.sock
chmod=0770
chown=${options.user}:${options.group}

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix://${parts.buildout.directory}/var/supervisor.sock ; use a unix:// URL  for a unix socket

[eventlistener:memmon]
command=/usr/local/bin/memmon -a 1024MB
events=TICK_60

[program:${options.domain}-nginx]
command = ${parts.buildout.directory}/parts/nginx/sbin/nginx
autostart = true
redirect_stderr = true
stdout_logfile = ${parts.buildout.directory}/var/log/nginx-stdout.log

[program:${options.domain}-zeo]
command=${parts.buildout.directory}/bin/zeo fg
autostart=true

{% for instance in options.instances.split() %}
[program:${options.domain}-client-${instance}]
command=${parts.buildout.directory}/bin/${instance} console
autostart=true
{% end %}

[program:haproxy]
command=${parts.buildout.directory}/bin/haproxy -f ${parts.buildout.directory}/etc/haproxy.conf -db
autostart=false
