# templates/supervisord.conf.in
[supervisord]
umask=002                   ; (process file creation umask;default 022)
logfile=/tmp/supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB       ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=10          ; (num of main logfile rotation backups;default 10)
loglevel=info               ; (logging level;default info; others: debug,warn)
pidfile=/tmp/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false              ; (start in foreground if true;default false)
minfds=1024                 ; (min. avail startup file descriptors;default 1024)
minprocs=200                ; (min. avail process descriptors;default 200)
user=${user}
enironment=PATH=/usr/local/sbin:/usr/local/bin:/usr/bin:/usr/sbin:/bin:/sbin:/usr/X11R6/bin
identifier=supervisor

[unix_http_server]
file=/tmp/supervisor.sock
chmod=0770
chown=${user}:${user}

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL  for a unix socket

[eventlistener:memmon]
command=/usr/local/bin/memmon -a 1024MB
events=TICK_60

[program:plone.org-nginx]
command = ${nginx:location}/sbin/nginx
autostart = true
redirect_stderr = true
stdout_logfile = ${buildout:directory}/var/log/nginx-stdout.log
priority = 3

[program:plone.org-zss]
command=/srv/plone.org/parts/zeo/bin/runzeo
priority=1
autostart=true

[program:plone.org-client-1]
command=/srv/plone.org/parts/instance1/bin/runzope
priority=2
autostart=true

[program:plone.org-client-2]
command=/srv/plone.org/parts/instance2/bin/runzope
priority=2
autostart=true

[program:plone.org-client-3]
command=/srv/plone.org/parts/instance3/bin/runzope
priority=2
autostart=true

[program:plone.org-client-4]
command=/srv/plone.org/parts/instance4/bin/runzope
priority=2
autostart=true

[program:plone.org-client-5]
command=/srv/plone.org/parts/instance5/bin/runzope
priority=2
autostart=false

[program:plone.org-client-6]
command=/srv/plone.org/parts/instance6/bin/runzope
priority=2
autostart=false

[program:plone.org-client-7]
command=/srv/plone.org/parts/instance7/bin/runzope
priority=2
autostart=false

[program:plone.org-client-8]
command=/srv/plone.org/parts/instance8/bin/runzope
priority=2
autostart=false
