# templates/supervisord.conf.in
[unix_http_server]
file = ${buildout:directory}/var/supervisor.sock
chmod = 0600
chown = root:admin

[supervisorctl]
serverurl = unix://${buildout:directory}/var/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisord]
logfile = ${buildout:directory}/var/log/supervisord.log
logfile_maxbytes = 5MB
logfile_backups = 10
loglevel = info
pidfile = ${buildout:directory}/var/supervisord.pid ;
childlogdir = ${buildout:directory}/var/log
nodaemon = false
directory = ${buildout:directory}/var/log

[program:nginx]
command = ${nginx:location}/sbin/nginx
environment = LD_LIBRARY_PATH=${nginx:libxml2}/lib:${nginx:libxslt}/lib
autostart = true
autorestart = true
stdout_logfile = NONE
stderr_logfile = NONE
priority = 40

[program:haproxy]
command = ${buildout:directory}/bin/haproxy -f ${buildout:directory}/etc/haproxy.conf -db
autostart= true
autorestart = true
redirect_stderr = true
priority = 30

[program:varnish]
command = ${varnish:location}/sbin/varnishd -f ${buildout:directory}/etc/varnish.vcl -F -a ${varnish-conf:listen} -T ${varnish-conf:telnet}
autostart= true
autorestart = true
redirect_stderr = true
priority = 5

[program:zeo]
command = ${zeo:location}/bin/runzeo 
directory = ${zeo:location}
priority = 10

[program:instance]
numprocs = 4
numprocs_start = 1
process_name = %(process_num)s
command = ${buildout:directory}/parts/instance%(process_num)s/bin/runzope 
priority = 20

[program:dev-instance]
command = ${dev-instance:location}/bin/runzope 
process_name = dev-instance
autostart = false
