# templates/supervisord-extras.conf.in
[program:haproxy]
command = ${buildout:directory}/bin/haproxy -f ${buildout:directory}/etc/haproxy.conf -db
autostart= false
autorestart = true
redirect_stderr = true
priority = 30

[program:varnish]
command = ${varnish:location}/sbin/varnishd -f ${buildout:directory}/etc/varnish.vcl -F -a ${varnish-conf:listen} -T ${varnish-conf:telnet}
autostart= false
autorestart = true
redirect_stderr = true
priority = 5
